# supervisor.conf
[supervisord]
nodaemon=true

[program:app]
command=make serve

[program:nginx]
command=nginx
